# files/systemd/physical-device-manager.service.in
[Unit]
Description=webos - "%n"
Requires=ls-hubd.service systemd-udevd.service
After=ls-hubd.service systemd-udevd.service

[Service]
Type=simple
OOMScoreAdjust=-500
EnvironmentFile=-/var/systemd/system/env/physical-device-manager.env
ExecStart=/usr/sbin/physical-device-manager
ExecStop=-/usr/bin/luna-send -n 1 luna://com.webos.service.pdm/umountAllDrive '{}'
Restart=on-failure
