# files/systemd/auto/physical-device-manager.service.in
[Unit]
Description=webos - "%n"
Requires=ls-hubd.service systemd-udevd.service
After=ls-hubd.service systemd-udevd.service

[Service]
Type=simple
OOMScoreAdjust=-500
EnvironmentFile=-/var/systemd/system/env/physical-device-manager.env
#ExecStartPre=/usr/sbin/setcap cap_sys_admin,cap_chown+eip /usr/sbin/physical-device-manager
ExecStart=/usr/sbin/physical-device-manager
ExecStop=-/usr/bin/luna-send -n 1 luna://com.webos.service.pdm/umountAllDrive '{}'
Restart=on-failure
User=pdmuser
Group=pdmgroup
#CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SETGID CAP_CHOWN CAP_DAC_OVERRIDE CAP_SETUID
AmbientCapabilities=CAP_SYS_ADMIN CAP_SETGID CAP_CHOWN CAP_DAC_OVERRIDE CAP_SETUID CAP_FOWNER
